# systemctl_services/roscore.service
# This file should be placed in: /etc/systemd/system/roscore.service
# Symlink: sudo ln -s /home/administrator/eeha/test_utils/systemctl_services/roscore.service /etc/systemd/system/

[Unit]
Description="ROS Master (roscore) - Central coordination service"
After=network.target
Wants=network.target

[Service]
Type=simple
User=administrator
Environment="ROS_MASTER_URI=http://localhost:11311"
ExecStart=/bin/bash -c "source /home/administrator/husky_setup.sh; _husky_check_setup; echo '[roscore.service] Starting ROS Master on port 11311'; /opt/ros/noetic/bin/roscore"
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
